# asd.cron.hourly
#!/bin/sh
/usr/bin/anything-sync-daemon sync &> /dev/null
