# asd.cron.hourly
#!/bin/sh
/usr/sbin/rc.d sync asd &> /dev/null
