# common/asd.cron.hourly
#!/bin/sh
/usr/bin/anything-sync-daemon resync &> /dev/null
